# dockerfiles/cron/app.crontab
* * * * * cd /app && PATH=/usr/local/bin:$PATH php vendor/bin/crunz schedule:run
# Empty line which is required apparently
